# deploy/systemd/trader-daily-cycle.timer
[Unit]
Description=Run trader daily cycle after market close

[Timer]
OnCalendar=Mon..Fri 16:05
Persistent=true
Unit=trader-daily-cycle.service

[Install]
WantedBy=timers.target
