# deploy/launchd/com.trader.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.trader.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/mananagarwal/Desktop/2nd brain/plant to image/trader/.venv/bin/python</string>
    <string>-m</string>
    <string>trader_system.runtime.local_api_server</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8050</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mananagarwal/Desktop/2nd brain/plant to image/trader</string>

  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mananagarwal/Desktop/2nd brain/plant to image/trader/cache/logs/launchd_dashboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mananagarwal/Desktop/2nd brain/plant to image/trader/cache/logs/launchd_dashboard.err.log</string>
</dict>
</plist>
